# supervisord.conf
[supervisord]
nodaemon=true

[program:kippo]
directory=/data/glastopf/
command=glastopf-runner
autorestart=true
