# dist/supervisord.conf
[supervisord]
nodaemon=true

[program:glastopf]
directory=/data/glastopf/
command=bash -c 'cp /opt/glastopf/glastopf.cfg /data/glastopf/ && exec glastopf-runner'
autorestart=true
